# umacore-bot.service
[Unit]
Description=UmaCore Discord Bot - Umamusume Club Quota Tracker
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/UmaCore
Environment="PATH=/home/pi/UmaCore/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/pi/UmaCore/.env

# Python virtual environment
ExecStart=/home/pi/UmaCore/venv/bin/python /home/pi/UmaCore/main.py

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=umacore-bot

# Security
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (optional, adjust as needed)
# LimitNOFILE=65536
# MemoryMax=512M

[Install]
WantedBy=multi-user.target

